# dp-switch/90-dp.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/etc/alsa/dp-switch.sh"
